# recipes-devtools/python/python3-labgrid/labgrid-exporter.service
[Unit]
Description=Labgrid exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/labgrid/environment
ExecStart=/usr/bin/python3 /usr/bin/labgrid-exporter -x ws://${LABGRID_COORDINATOR_IP}:${LABGRID_COORDINATOR_PORT}/ws /etc/labgrid/configuration.yaml

[Install]
WantedBy=multi-user.target
